# TASK.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2023-03-22T13:29:21.3125364</Date>
    <Author>{{COMPUTER HOSTNAME}}/{{YOUR USERNAME}}</Author>
    <Description>Sunshine setup requires specification of audio sync and virtual sync, but these change every reboot. This script will automatically update these values in Sunshine configuration for us.</Description>
    <URI>\Update Sunshine Audio Syncs</URI>
  </RegistrationInfo>
  <Triggers>
    <BootTrigger>
      <Enabled>true</Enabled>
    </BootTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-18</UserId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>StopExisting</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT1H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>{{ QUOTED PATH TO NODEJS npm.cmd (Typically "C:\Program Files\nodejs\npm.cmd") }}</Command>
      <Arguments>start set-audio-sink "{{ Full string name of audio sink device }}"</Arguments>
      <WorkingDirectory>{{ DIRECTORY TO THIS PROJECT, NOT IN QUOTES AND WITH A TRAILING SLASH}}</WorkingDirectory>
    </Exec>
    <Exec>
      <Command>{{ QUOTED PATH TO NODEJS npm.cmd (Typically "C:\Program Files\nodejs\npm.cmd") }}</Command>
      <Arguments>start set-virtual-sink "{{ Full string name of the virtual sink device }}"</Arguments>
      <WorkingDirectory>{{ DIRECTORY TO THIS PROJECT, NOT IN QUOTES AND WITH A TRAILING SLASH}}</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
